# payload/Library/LaunchAgents/com.erikng.umad.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.erikng.umad</string>
	<key>LimitLoadToSessionType</key>
	<array>
		<string>Aqua</string>
	</array>
	<key>ProgramArguments</key>
	<array>
		<string>/Library/umad/Resources/umad</string>
		<!-- Cutoff Date -->
		<string>--cutoffdate</string>
		<string>2023-10-16-17:00</string>
		<!-- Profile Identifier -->
		<!--<string>--profileidentifier</string> -->
		<!--<string>B68ABF1E-70E2-43B0-8300-AE65F9AFA330</string> -->
		<!-- -->
		<!-- Logos -->
		<string>--logopath</string>
		<string>/usr/local/nyuad/nyuadnudge.png</string>
		<string>--nagsspath</string>
		<string>/usr/local/nyuad/nag_ss.png</string>
		<string>--uasspath</string>
		<string>/usr/local/nyuad/uamdm.png</string>
		<!-- -->
		<!-- Text-->
		<!--<string>--depfailuretext</string> -->
		<!--<string>Not getting this notification?</string> -->
		<string>--depfailuresubtext</string>
		<string>Click the system clock to check or reboot mac or do manual enroll.</string>
		<string>--manualenrollmenttext</string>
		<string>Manual Enrollment Required - Click Button below and use "nyuad" for the GROUP ID</string>
		<string>--manualenrollh1text</string>
		<string>For Support</string>
		<string>--manualenrollh2text</string>
		<string>Email us nyuad.it@nyu.edu or visit the IT Support Center in C2</string>
		<string>--titletext</string>
		<string>NYUAD IT - Device Enrollment</string>
		<string>--subtitletext</string>
		<string>As per the email you received from NYUAD IT and the Intranet post dated 5th October 2023 titled “Computer Management Software - Action Required”</string>
		<string>--duedatetext</string>
		<string>You are required to enroll your computer in the new computer management tool to ensure compliance</string>
		<string>--paragraph1</string>
		<string>You can click on the “Guide Me” button in the bottom left of this window at any time for step-by-step instructions</string>
		<string>--paragraph2</string>
		<string>This process only takes 5 minutes and should be completed by the due date mentioned on the left. No restart of your computer is required.</string>
		<string>--paragraph3</string>
		<string>To get started, select the notification located at the upper-right corner of your screen and enter your Mac login credentials, or click the "Manual Enrollment" button below if directed to do so</string>
		<!--<string>--sysprefsh1text</string> -->
		<!--<string>Want this box to go away?</string> -->
		<!--<string>--sysprefsh2text</string> -->
		<!--<string>Open System Preferences and approve Device Management.</string> -->
		<!--<string>--uamdmparagraph1</string> -->
		<!--<string>Thank you for enrolling your device into MDM. We sincerely appreciate you doing this in a timely manner.</string> -->
		<!--<string>--uamdmparagraph2</string> -->
		<!--<string>Unfortunately, your device has been detected as only partially enrolled into our system.</string> -->
		<!--<string>--uamdmparagraph3</string> -->
		<!--<string>Please go to System Preferences -> Profiles, click on the Device Enrollment profile and click on the approve button.</string> -->
		<!-- -->
		<!-- URLs -->
		<string>--manualenrollmenturl</string>
		<string>https://ds1688.awmdm.com/DeviceManagement/Enrollment</string>
		<string>--moreinfourl</string>
		<string>https://docs.google.com/document/d/1_AKSNSMv_j23ek88ldsxnGwLeOBHxmb35DNcbQPvYyk/edit?usp=sharing</string>
		<!-- -->
		<!-- Timers -->
		<!--<string>--cutoffdatewarning</string> -->
		<!--<string>3</string> -->
		<!--<string>--timerday1</string> -->
		<!--<string>600</string> -->
		<!--<string>--timerday3</string> -->
		<!--<string>7200</string> -->
		<string>--timerelapsed</string>
		<string>600</string>
		<string>--timerfinal</string>
		<string>600</string>
		<!--<string>--timerinital</string> -->
		<!--<string>14400</string> -->
		<!--<string>--timermdm</string> -->
		<!--<string>5</string> -->
		<!-- -->
		<!-- Action_True -->
		<!--<string>--enableenrollmentbutton</string> -->
		<!--<string>--notimer</string> -->
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardOutPath</key>
	<string>/Library/umad/Logs/umad.log</string>
	<key>StandardErrorPath</key>
	<string>/Library/umad/Logs/umad.log</string>
	<key>StartCalendarInterval</key>
	<array>
		<dict>
			<key>Minute</key>
			<integer>0</integer>
		</dict>
		<dict>
			<key>Minute</key>
			<integer>30</integer>
		</dict>
	</array>
</dict>
</plist>

