# payload/Library/LaunchAgents/com.erikng.umad.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.erikng.umad</string>
	<key>LimitLoadToSessionType</key>
	<array>
		<string>Aqua</string>
	</array>
	<key>ProgramArguments</key>
	<array>
		<string>/Library/umad/Resources/umad</string>
		<!-- Cutoff Date -->
		<string>--cutoffdate</string>
		<string>2023-09-13-17:00</string>
		<!-- Profile Identifier -->
		<!--<string>--profileidentifier</string> -->
		<!--<string>B68ABF1E-70E2-43B0-8300-AE65F9AFA330</string> -->
		<!-- -->
		<!-- Logos -->
		<string>--logopath</string>
		<string>/usr/local/nyuad/nyuadnudge.png</string>
		<string>--nagsspath</string>
		<string>/usr/local/nyuad/nag_ss.png</string>
		<string>--uasspath</string>
		<string>/usr/local/nyuad/uamdm.png</string>
		<!-- -->
		<!-- Text-->
		<!--<string>--depfailuretext</string> -->
		<!--<string>Not getting this notification?</string> -->
		<!--<string>--depfailuresubtext</string> -->
		<!--<string>You can also enroll manually below:</string> -->
		<!--<string>--manualenrollmenttext</string> -->
		<!--<string>Manual Enrollment Required</string> -->
		<!--<string>--manualenrollh1text</string> -->
		<!--<string>Want this box to go away?</string> -->
		<string>--manualenrollh2text</string>
		<string>Click Manual Enrollment button below and use GROUP ID = aditmx.</string>
		<string>--titletext</string>
		<string>NYUAD IT - MDM Enrollment</string>
		<string>--subtitletext</string>
		<string>As communicated in the email titled "Computer Management Software - Action Required".</string>
		<string>--duedatetext</string>
		<string>Enrollment into MDM is mandatory to ensure that NYUAD IT can protect your computer with security software.</string>
		<string>--paragraph1</string>
		<string>Please complete the action before the due date and it takes only 5 Minutes of your time. No restart is required.</string>
		<string>--paragraph2</string>
		<string>To enroll, click on the notification in the top right corner similar as per the sample image below, then click allow and provide your mac login password.</string>
		<string>--paragraph3</string>
		<string>Click "More Info" for guide and please contact nyuad.it@nyu.edu for issues. Thank You</string>
		<!--<string>--sysprefsh1text</string> -->
		<!--<string>Want this box to go away?</string> -->
		<!--<string>--sysprefsh2text</string> -->
		<!--<string>Open System Preferences and approve Device Management.</string> -->
		<!--<string>--uamdmparagraph1</string> -->
		<!--<string>Thank you for enrolling your device into MDM. We sincerely appreciate you doing this in a timely manner.</string> -->
		<!--<string>--uamdmparagraph2</string> -->
		<!--<string>Unfortunately, your device has been detected as only partially enrolled into our system.</string> -->
		<!--<string>--uamdmparagraph3</string> -->
		<!--<string>Please go to System Preferences -> Profiles, click on the Device Enrollment profile and click on the approve button.</string> -->
		<!-- -->
		<!-- URLs -->
		<string>--manualenrollmenturl</string>
		<string>https://ds1688.awmdm.com/DeviceManagement/Enrollment</string>
		<string>--moreinfourl</string>
		<string>https://docs.google.com/document/d/1_AKSNSMv_j23ek88ldsxnGwLeOBHxmb35DNcbQPvYyk/edit?usp=sharing</string>
		<!-- -->
		<!-- Timers -->
		<string>--cutoffdatewarning</string>
		<string>3</string>
		<string>--timerday1</string>
		<string>28</string>
		<!--<string>--timerday3</string> -->
		<!--<string>7200</string> -->
		<!--<string>--timerelapsed</string> -->
		<!--<string>10</string> -->
		<!--<string>--timerfinal</string> -->
		<!--<string>60</string> -->
		<!--<string>--timerinital</string> -->
		<!--<string>14400</string> -->
		<!--<string>--timermdm</string> -->
		<!--<string>5</string> -->
		<!-- -->
		<!-- Action_True -->
		<!--<string>--enableenrollmentbutton</string> -->
		<!--<string>--notimer</string> -->
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardOutPath</key>
	<string>/Library/umad/Logs/umad.log</string>
	<key>StandardErrorPath</key>
	<string>/Library/umad/Logs/umad.log</string>
	<key>StartCalendarInterval</key>
	<array>
		<dict>
			<key>Minute</key>
			<integer>0</integer>
		</dict>
		<dict>
			<key>Minute</key>
			<integer>30</integer>
		</dict>
	</array>
</dict>
</plist>

